# deployment/crypto-trading-bot-5x3pct.service
[Unit]
Description=Crypto Trading Bot 5x 3pct
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/projects/crypto-trading-project
ExecStart=/usr/local/bin/python3 /root/projects/crypto-trading-project/execution/run_bot.py --config /root/projects/crypto-trading-project/execution/config.live.5x-3pct.json --run-loop
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/root/projects/crypto-trading-project/live_bot_5x3pct.log
StandardError=append:/root/projects/crypto-trading-project/live_bot_5x3pct.log

[Install]
WantedBy=multi-user.target
